# docker/supervisord.conf
[supervisord]
nodaemon=true
loglevel=info

[group:litellm]
programs=main,health

[program:main]
command=sh -c 'if [ "$USE_DDTRACE" = "true" ]; then export DD_TRACE_OPENAI_ENABLED="False"; exec ddtrace-run uvicorn litellm.proxy.proxy_cli:main --host 0.0.0.0 --port=4000; else exec uvicorn litellm.proxy.proxy_cli:main --host 0.0.0.0 --port=4000; fi'
autostart=true
autorestart=false
startretries=0
priority=1
exitcodes=0
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
environment=PYTHONUNBUFFERED=true

[program:health]
command=sh -c '[ "$SEPARATE_HEALTH_APP" = "1" ] && exec uvicorn litellm.proxy.health_endpoints.health_app_factory:build_health_app --factory --host 0.0.0.0 --port=${SEPARATE_HEALTH_PORT:-4001} || exit 0'
autostart=true
autorestart=false
startretries=0
priority=2
exitcodes=0
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
environment=PYTHONUNBUFFERED=true